# ansible/roles/monitoring/templates/observer-exporter.service.j2
[Unit]
Description=P2Pool Observer Prometheus Exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=exporter
Group=exporter
WorkingDirectory={{ observer_exporter_dir }}
ExecStart=/usr/bin/python3 {{ observer_exporter_dir }}/p2pool_observer_exporter.py \
  --listen-address {{ observer_exporter_bind | default('127.0.0.1') }} \
  --listen-port {{ observer_exporter_port }} \
  --log-level INFO
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# plan-approved: ansible-role-edit -- H-001 wallet leak fix (mr-hardass 2026-05-18):
# wallet must NOT be inlined here — systemctl show / systemctl cat /
# journalctl unit dumps render the .service file world-readable.
# EnvironmentFile= is loaded by systemd PID 1 (mode 0600 root:root OK);
# downstream User=exporter still receives the env vars after privilege drop.
EnvironmentFile=/etc/monero-farm/p2pool-wallet.env
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
