# ansible/roles/monitoring/templates/observer-exporter.service.j2
[Unit]
Description=P2Pool Observer Prometheus Exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=exporter
Group=exporter
WorkingDirectory={{ observer_exporter_dir }}
ExecStart=/usr/bin/python3 {{ observer_exporter_dir }}/p2pool_observer_exporter.py \
  --listen-address {{ observer_exporter_bind | default('127.0.0.1') }} \
  --listen-port {{ observer_exporter_port }} \
  --log-level INFO
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=MINER_ADDRESS={{ monero_wallet_address }}
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
