# ansible/roles/monitoring/templates/unified-exporter.service.j2
[Unit]
Description=Unified XMR Exporter (monerod + p2pool main + p2pool mini)
After=network.target

[Service]
User=exporter
Group=exporter
ExecStart=/usr/bin/python3 {{ unified_exporter_dir }}/p2pool_exporter.py \
  --listen-address {{ unified_exporter_bind }} \
  --listen-port {{ unified_exporter_port }}
Restart=always
RestartSec=3
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
